# .services/runit/openvpn-client/finish
#!/bin/bash

killall -SIGINT openvpn
rm -f /var/log/openvpn.log
