# .services/runit/killswitch/finish
#!/bin/bash

exec /bin/ovpn --dev-function Killswitch_Disable
killall -SIGINT openvpn
